# Firmware - KKEPS-3.5 Extracted/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
sleep 5
/sbin/kkeps_on &
/sbin/kkeps_off &
/sbin/kkeps_seekwifi &
/sbin/kkeps_reconect &

exit 0
